# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Run database migrations on startup (MUST succeed before services start)
[program:init-database]
command=/app/scripts/init-database.sh
directory=/app
user=root
autostart=true
autorestart=false
startsecs=0
exitcodes=0
priority=10
redirect_stderr=true
stdout_logfile=/var/log/supervisor/init-database.log
stderr_logfile=/var/log/supervisor/init-database.error.log
environment=UV_THREADPOOL_SIZE=64

# MCP Server (Core workflow engine)
[program:mcp-server]
command=/bin/sh -c "/app/scripts/wait-for-init.sh && exec /usr/local/bin/tsx packages/mcp-server/src/server.ts"
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/mcp-server.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2
environment=NODE_ENV=production,METRICS_PORT=9090,UV_THREADPOOL_SIZE=64,NODE_OPTIONS="--max-old-space-size=512",LOG_FILE=/var/log/app/mcp-server.log
priority=50

# Backend API Server
[program:backend-api]
command=/bin/sh -c "/app/scripts/wait-for-init.sh && exec /usr/local/bin/tsx packages/web-backend/src/server.ts"
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/backend-api.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2
environment=NODE_ENV=production,METRICS_PORT=9090,UV_THREADPOOL_SIZE=64,NODE_OPTIONS="--max-old-space-size=512",LOG_FILE=/var/log/app/backend-api.log
priority=50

# Nginx for frontend static files
[program:nginx-frontend]
command=/bin/sh -c "/app/scripts/wait-for-init.sh && exec /usr/sbin/nginx -g 'daemon off;'"
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.error.log
priority=50

# Startup health check (runs once after all services started)
[program:startup-health-check]
command=/bin/sh -c "/app/scripts/wait-for-init.sh && sleep 5 && /app/scripts/health-check.sh && touch /tmp/app-ready"
directory=/app
user=root
autostart=true
autorestart=false
startsecs=0
exitcodes=0
priority=100
redirect_stderr=true
stdout_logfile=/var/log/supervisor/startup-health.log
stderr_logfile=/var/log/supervisor/startup-health.error.log

# Group all services for easy management
[group:moira-services]
programs=init-database,mcp-server,backend-api,nginx-frontend,startup-health-check
priority=999